# seafile.conf
description     "start seafile"
start on startup
task
exec /seafile/seafile-server-latest/seafile.sh start
