# seafile.conf
description     "my job"
start on startup
task
exec /path/to/my/script.sh
