# seahub.conf
description     "start seahub"
start on startup
task
exec /seafile/seafile-server-latest/seahub.sh start
